# templates/default/delegate.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:        iproute2
# Required-Start:    
# Required-Stop:     
# Default-Start:   2 3 4 5
# Default-Stop:    
# Short-Description: Start DELEGATE daemon
### END INIT INFO

# *****************************************************************
# Dynamic iproute2 configuration  --  DO NOT EDIT THIS FILE BY HAND
# Generated by chef - <%= @date %>
# *****************************************************************
PATH=/sbin:/bin:/usr/sbin:/usr/bin

. /lib/lsb/init-functions

DAEMON=/usr/local/sbin/delegated
PIDFILE=/var/run/delegate.pid

test -x $DAEMON || exit 5

DELEGATE_OPTS=<% @options.each do |i| %><%= i %> <% end -%>

case "$1" in
 start)
  log_daemon_msg "Starting DELEGATE daemon" "delegate"
  start-stop-daemon --start --quiet --oknodo -m --pidfile $PIDFILE --startas $DAEMON -- $DELEGATE_OPTS
  log_end_msg $?
  ;;
 stop)
  log_daemon_msg "Stopping DELEGATE daemon" "delegate"
  start-stop-daemon --stop --quiet --oknodo --pidfile $PIDFILE
  log_end_msg $?
  rm -f $PIDFILE
  ;;
 restart|force-reload)
  $0 stop && sleep 2 && $0 start
  ;;
 status)
  status_of_proc $DAEMON "delegate daemon"
  ;;
 *)
  echo "Syntaxe: $(basename $0) [start|stop|status]..." && false
esac

